# tools/pixelplein-kiosk.service
[Unit]
Description=PixelPlein Kiosk Browser
After=network.target pixelplein-provisioner.service
Wants=pixelplein-provisioner.service

[Service]
User=pixelplein
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pixelplein/.Xauthority
ExecStart=/opt/pixelplein/tools/pixelplein-kiosk.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
